# l850gl/files/etc/hotplug.d/usb/99-l850gl
#!/bin/sh
#
# Hotplug handler untuk Fibocom L850-GL (dalam mode NCM)
# VID:PID = 8087:095a
#

case "$ACTION" in
    add|bind)
        case "$PRODUCT" in
            8087/095a/*)
                logger -t "l850gl" "Modem L850GL (mode NCM) terdeteksi (PRODUCT=$PRODUCT), menjalankan inisialisasi..."
                /usr/bin/l850gl-init add &
                ;;
        esac
        ;;
    remove)
        case "$PRODUCT" in
            8087/095a/*)
                logger -t "l850gl" "Modem L850GL (mode NCM) dicabut (PRODUCT=$PRODUCT), membersihkan interface..."
                /usr/bin/l850gl-init remove &
                ;;
        esac
        ;;
esac

exit 0
