# l850gl/files/etc/init.d/l850gl-init
#!/bin/sh /etc/rc.common

START=99
USE_PROCD=1

boot() {
    logger -t "l850gl" "Init script: menjalankan l850gl-init add setelah boot"
    (sleep 15 && /usr/bin/l850gl-init add) &
}
